# Windows/Suspicious ScreenSaver Registry Key Set Detected.json
[
  {
    "name": "Suspicious ScreenSaver Registry Key Set Detected",
    "description": "Detection of suspicious modification to the SCRNSAVE.EXE registry key, indicating that a file is being set as the system screensaver, often for persistence or execution.",
    "tags": [
      "security",
      "detection",
      "medium",
      "windows",
      "persistence",
      "defense_evasion",
      "screensaver",
      "registry",
      "lolbin",
      "TA0005",
      "TA0003",
      "T1547",
      "T1547.001",
      "T1218",
      "T1218.011"
    ],
    "prompts": [
      {
        "id": "1",
        "prompt": "MEDIUM SEVERITY ALERT: Suspicious ScreenSaver Registry Key Set Detected.\\n\\nThe `SCRNSAVE.EXE` registry key was suspiciously modified by `rundll32.exe` on host '{ComputerName}' by user '{User}' at {Timestamp}. The key '{TargetObject}' was set to '{Details}', indicating a new screen saver file has been configured. This is a known technique for persistence and defense evasion.\\n\\n**Event Details:**\\n* **Timestamp**: {Timestamp}\\n* **Host**: {ComputerName}\\n* **User**: {User}\\n* **Process**: {Image}\\n* **Target Registry Key**: {TargetObject}\\n* **New Value (Screensaver Path)**: {Details}\\n* **Sysmon Event ID**: {EventID}\\n* **Sysmon Rule Name**: {RuleName}\\n\\n**RECOMMENDED IMMEDIATE ACTIONS (PRIORITY MEDIUM):**\\n\\n1.  **VERIFY LEGITIMACY:** Determine if this screensaver change was a legitimate action by the user '{User}' or an authorized system administrator/management tool (e.g., GPO, desktop management software). Contact the user directly for confirmation. \\n2.  **ANALYZE SCREENSAVER FILE:** If the screensaver path '{Details}' is not a known legitimate system file (i.e., not in `C:\\Windows\\System32\\` or `C:\\Windows\\SysWOW64\\`), immediately obtain a copy of the file for analysis. Perform static and dynamic analysis in a sandbox environment (e.g., submit to VirusTotal). Look for malicious indicators, obfuscation, or unusual behavior.\\n3.  **INVESTIGATE USER ACTIVITY:** Investigate the '{User}' account's recent activity on '{ComputerName}' for other suspicious actions leading up to this event (e.g., execution of unknown files, unusual network connections, phishing attempts).\\n4.  **CHECK FOR PERSISTENCE:** Even if the file itself is not immediately malicious, its placement as a screensaver could indicate an attempt at persistence. Check startup locations, scheduled tasks, and other common persistence mechanisms on '{ComputerName}'.\\n5.  **REVIEW GROUP POLICIES/MDM:** If applicable, verify whether any Group Policy Objects (GPOs) or Mobile Device Management (MDM) solutions (like Intune) are configured to manage screensaver settings. Deviations could indicate malicious activity.\\n6.  **FORENSIC ANALYSIS (IF MALICIOUS):** If the screensaver file or associated activity is confirmed malicious, proceed with a full forensic investigation of '{ComputerName}' to identify the initial compromise vector, scope of impact, and other IOCs.\\n7.  **THREAT HUNT:** Search your environment for other instances of the `SCRNSAVE.EXE` registry key being modified, especially if `rundll32.exe` is the modifying process and the path is outside standard system directories.",
        "continueOnError": false,
        "systemCapability": "analyzeSecurityEvent"
      }
    ],
    "parameters": [
      {
        "name": "Timestamp",
        "required": true,
        "description": "The timestamp of when the registry key was modified."
      },
      {
        "name": "Image",
        "required": true,
        "description": "The process (expected to be rundll32.exe) that modified the registry key."
      },
      {
        "name": "TargetObject",
        "required": true,
        "description": "The full path to the modified registry key."
      },
      {
        "name": "Details",
        "required": true,
        "description": "The new value set for the registry key (the path to the .scr file)."
      },
      {
        "name": "User",
        "required": true,
        "description": "The user account under which the registry modification occurred."
      },
      {
        "name": "ComputerName",
        "required": true,
        "description": "The hostname of the system where the event occurred."
      },
      {
        "name": "EventID",
        "required": true,
        "description": "The Sysmon event ID (expected to be 13)."
      },
      {
        "name": "RuleName",
        "required": false,
        "description": "The Sysmon rule name that triggered the alert, if available."
      }
    ]
  }
]