# Database/TA0005_Persistencia/Suspicious Scheduled Task Creation via Masqueraded XML.yml
title: Suspicious Scheduled Task Creation via Masqueraded XML
id: 5e6f7a8b-9c0d-4e1f-2a3b-4c5d6e7f8a9b # Replace with a newly generated UUID
description: >
  The creation of a suspicious scheduled task using an XML file with a masqueraded extension is detected.
  This indicates an attempt to establish persistence or execute code covertly.
status: experimental
author: aguztynbassi
date: 2025/04/05
modified: 2025/05/25 # Updated to current date
logsource:
  category: process_creation # Standard Sigma category for process creation events
  product: windows
  service: sysmon # Specifies Sysmon as the logging service
detection:
  selection:
    EventID: 1 # Sysmon Event ID for Process Creation
    Image|endswith: '\schtasks.exe' # Use |endswith for precise binary name
    CommandLine|contains: # Looking for creation command with XML flag
      - '/create'
      - '-create'
    CommandLine|contains|all: # Ensuring both create and XML flags are present
      - '/xml'
      - '-xml'
  # Filters for legitimate scheduled task creations
  # It's important to refine these filters based on your environment to minimize false positives.
  filter_system_integrity:
    IntegrityLevel: 'System' # Exclude tasks created by System integrity level, often legitimate
  filter_command_xml_extension:
    CommandLine|endswith: '.xml"' # Exclude if the XML file explicitly ends with .xml"
    # This filter specifically targets the "masqueraded" part: an XML file
    # that does NOT explicitly have a .xml extension in the command line,
    # or where the extension is hidden/masqueraded.
    # The original "filter3" was `command: "*.xml*"`, which would actually
    # *include* legitimate XML files. This filter aims to *exclude* them.
  filter_legitimate_parents:
    ParentImage|startswith: # Common legitimate parent processes for scheduled tasks
      - 'C:\ProgramData\OEM\UpgradeTool\CareCenter_'
      - 'C:\Program Files\Axis Communications\AXIS Camera Station\'
      - 'C:\Program Files\Axis Communications\AXIS Device Manager\'
      - 'C:\Program Files (x86)\Zemana\AntiMalware\'
      - 'C:\Program Files\Dell\SupportAssist\'
    ParentImage|endswith:
      - '\Setup_msi.exe'
      - '\SetupActions.exe'
      - '\AdmSetupActions.exe'
      - '\AntiMalware.exe'
      - '\pcdrcui.exe'
  filter_rundll32_installer:
    ParentImage|endswith: '\rundll32.exe'
    CommandLine|contains: ':\WINDOWS\Installer\MSI'
    CommandLine|endswith: ',zzzzInvokeManagedCustomActionOutOfProc'
  condition: selection and not (filter_system_integrity or filter_command_xml_extension or filter_legitimate_parents or filter_rundll32_installer)
fields:
  - Image
  - ParentImage
  - CommandLine
  - User
  - IntegrityLevel
  - RuleName # Sysmon specific field for rule matching
falsepositives:
  - Legitimate software installations or updates that use schtasks with XML for unusual reasons.
  - Custom internal scripts or deployment tools that create scheduled tasks with non-standard XML file naming conventions.
level: high
tags:
  - attack.tactic: TA0003 # Persistence
  - attack.tactic: TA0005 # Defense Evasion
  - attack.technique: T1053 # Scheduled Task/Job
  - attack.technique: T1053.005 # Scheduled Task/Job: Scheduled Task
  - attack.technique: T1036 # Masquerading
references:
  - https://attack.mitre.org/techniques/T1053/005/
  - https://docs.microsoft.com/en-us/windows/security/threat-protection/auditing/event-4698 # While this rule uses Sysmon, 4698 is related to scheduled task creation
  - https://learn.microsoft.com/en-us/sysinternals/downloads/sysmon # Sysmon official page