# Database/TA0003_Acceso_Inicial/PowerShell Obfuscation.yml
title: Registry Run Key Persistence
id: 5a8e0f9b-6b2c-4e3a-9c7d-1f2e3d4c5b6a # Generated UUID
description: >
  Detection has identified changes to the Registry Run keys for persistence.
  Adversaries may use this technique to maintain access by executing malicious binaries during system startup.
status: experimental
author: aguztynbassi
date: 2025/04/05
modified: 2025/04/05
logsource:
  product: windows
  service: security
detection:
  selection:
    EventID: 4663
    ObjectName|contains: '\REGISTRY\MACHINE\SOFTWARE\Microsoft\Windows\CurrentVersion\Run'
    # EventID 4663 is "An attempt was made to access an object."
    # We are looking for writes to the Run keys. While 4663 shows access,
    # the specific AccessMask for write operations should ideally be included
    # for more precise detection. Common write masks include:
    # 0x2 (FILE_WRITE_DATA)
    # 0x10000 (WRITE_DAC)
    # 0x20000 (WRITE_OWNER)
    # 0x20006 (KEY_SET_VALUE) - for registry writes specifically
    # For a more robust rule, consider adding:
    # AccessMask:
    #   - '0x2'
    #   - '0x10000'
    #   - '0x20000'
    #   - '0x20006'
  condition: selection
fields:
  - SubjectUserName
  - ProcessName
  - ObjectName
  - AccessMask
falsepositives:
  - Legitimate software installations
  - System updates
level: high
tags:
  - attack.tactic: TA0003 # Persistence
  - attack.technique: T1547.001 # Registry Run Keys / Startup Folder
references:
  - https://docs.microsoft.com/en-us/windows/security/threat-protection/auditing/event-4663
  - https://attack.mitre.org/techniques/T1547/001/