# Database/TA0003_Acceso_Inicial/Scheduled Task Creation.yml
title: Scheduled Task Creation
id: 1e8e2c7b-6f2a-4c5b-9c1a-2d3e4f5a6b7c  # Reemplaza con un UUID generado si lo deseas
description: Detects the creation of new scheduled tasks which may indicate persistence techniques.
status: experimental
author: aguztynbassi
date: 2024/05/25
logsource:
  product: windows
  service: security
detection:
  selection:
    EventID: 4698
  condition: selection
fields:
  - SubjectUserName
  - TaskName
  - TaskContent
falsepositives:
  - Legitimate scheduled task creation by administrators or system processes
level: medium
tags:
  - attack.tactic: TA0003 # Persistence
  - attack.technique: T1053 # Scheduled Task/Job
  - attack.technique: T1053.005 # Scheduled Task/Job: Scheduled Task Creation
references:
  - https://attack.mitre.org/techniques/T1053/ # MITRE ATT&CK T1053 (Scheduled Task/Job)
  - https://learn.microsoft.com/en-us/windows/security/threat-protection/auditing/event-4698 # Windows Event ID 4698 documentation
#   - https://www.fireeye.com/blog/threat-research/2020/04/abusing-scheduled-tasks-for-persistence.html # FireEye blog on scheduled tasks
  - https://www.redcanary.com/blog/scheduled-task-persistence/ # Red Canary blog on scheduled task persistence